# systemd/kanshi.service
[Unit]
Description=Kanshi - Screen Output Handler

[Service]
Type=simple
ExecStart=/usr/bin/kanshi

[Install]
WantedBy=river.target
